# laptop/.config/systemd/user/redshift.service
[Unit]
Description=Redshift Daemon

[Service]
ExecStart=/usr/bin/redshift

[Install]
WantedBy=default.target
